# boot/crontab
@reboot sh /home/pi/EVO-SmartIoTVehicle/launcher.sh